# system_files/usr/lib/systemd/system/bootc-fetch-apply-updates.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=*-*-* 12:00:00
